# ChatRoom/flask/flask.service
[Unit]
Description=Example systemd service.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/flask.sh

[Install]
WantedBy=multi-user.target
